# iso/airootfs/etc/systemd/system/lytux-install.service
[Unit]
Description=Run Lytux Installer After Network Connectivity
After=multi-user.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/lytux-install

[Install]
WantedBy=multi-user.target
